# dot_config/systemd/remove_user/handy.service
[Unit]
Description=Handy STT transcription daemon
After=graphical-session.target

[Service]
Type=exec
ExecStart=handy --start-hidden
Restart=on-failure
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
